# openrc/init.d/network
#!/sbin/runscript
# Copyright (c) 2009 Roy Marples <roy@marples.name>
# All rights reserved. Released under the 2-clause BSD license.

# This script was inspired by the equivalent rc.d network from NetBSD.

description="Configures network interfaces."

depend()
{
  #need localmount
  after rsyslog
  provide net
  keyword -jail -prefix -vserver
}

start()
{
  # bring up the loopback interface
  if [ -d /sys/class/net/lo ]; then
      ebegin "Bringing up loopback interface"
        /sbin/ifconfig lo 127.0.0.1 up
        if [ $? -ne 0 ]; then
            eend $?
        else
            eend 0
        fi
  fi
}
